# deploy/sora-missioncontrol-proxy.service
[Unit]
Description=Sora Mission Control same-origin app and admin proxy
After=network-online.target
Wants=network-online.target

# ── Deployment notes ──────────────────────────────────────────────────
#
# This service file is for PRODUCTION deployment on baset-ai (192.168.0.85).
# For Unraid (PRIMARY production), use Docker + docker-compose.yml instead.
# For development on baset-ai, use deploy/sora-missioncontrol-proxy-dev.service
# which runs as a user service with optional auth mode.
#
# Required env vars (in EnvironmentFile or Environment):
#   MISSION_CONTROL_ADMIN_PROXY_KEY  — admin proxy token (keep secret)
#
# Optional env vars (defaults shown):
#   MISSION_CONTROL_PROXY_HOST=0.0.0.0
#   MISSION_CONTROL_PROXY_PORT=3187
#   MISSION_CONTROL_PROXY_AUTH_MODE=required
#   HERMES_DASHBOARD_URL=http://192.168.0.85:9119
#   HERMES_WS_URL=ws://192.168.0.85:9119
#   MISSION_CONTROL_CORS_ORIGIN  — unset (same-host/loopback CORS)
#
# Unraid Docker env vars (set in docker-compose.yml, not here):
#   HERMES_REMOTE_USER, HERMES_REMOTE_HOST, HERMES_REMOTE_PATH,
#   HERMES_SSH_KEY, HERMES_DASHBOARD_PROXY_TARGET, HERMES_ACCEPT_HOOKS
#
# ──────────────────────────────────────────────────────────────────────

[Service]
Type=simple
WorkingDirectory=/opt/sora-missioncontrol/current
Environment=NODE_ENV=production
Environment=MISSION_CONTROL_PROXY_HOST=0.0.0.0
Environment=MISSION_CONTROL_PROXY_PORT=3187
Environment=MISSION_CONTROL_PROXY_AUTH_MODE=required
Environment=HERMES_DASHBOARD_URL=http://192.168.0.85:9119
Environment=HERMES_WS_URL=ws://192.168.0.85:9119
# Runtime secrets (MISSION_CONTROL_ADMIN_PROXY_KEY) from operator-managed file
EnvironmentFile=-/etc/sora-missioncontrol/proxy.env
ExecStart=/home/wliob/.hermes/node/bin/node missionControlProxy.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
